# src/main/assembly/dist/install/dd-data-vault.service
[Unit]
Description=Data Vault Service

[Service]
ExecStart=/opt/dans.knaw.nl/dd-data-vault/bin/dd-data-vault server /etc/opt/dans.knaw.nl/dd-data-vault/config.yml
# Java returns 143 even if the SIGTERM was handled correctly.
SuccessExitStatus=143
Environment="INITIAL_HEAP_SIZE=-Xms64m"
Environment="LC_ALL=en_US.UTF-8"

User=dd-data-vault
Group=dd-data-vault

[Install]
WantedBy=multi-user.target
